# roles/host_networking/templates/routerlab-network.service.j2
# Managed by Ansible (host_networking role)
[Unit]
Description=RouterLab host network preparation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ routerlab_script_path }}
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
